# package/easymesh-controller/files/etc/uci-defaults/90-easymesh-controller-enable
#!/bin/sh

uci -q get easymesh.main.enabled >/dev/null 2>&1 || uci set easymesh.main.enabled='1'
uci -q get easymesh.main.pairing_enabled >/dev/null 2>&1 || uci set easymesh.main.pairing_enabled='1'
uci -q get easymesh.main.local_member >/dev/null 2>&1 || uci set easymesh.main.local_member='0'
uci -q get easymesh.main.network_mode >/dev/null 2>&1 || uci set easymesh.main.network_mode='gateway'
uci commit easymesh >/dev/null 2>&1 || true

/etc/init.d/easymesh-controller enable >/dev/null 2>&1 || true
/etc/init.d/uhttpd enable >/dev/null 2>&1 || true
/etc/init.d/uhttpd restart >/dev/null 2>&1 || true

# Advertise the AC over mDNS so AP agents can discover it without a hardcoded URL.
/etc/init.d/umdns enable >/dev/null 2>&1 || true
/etc/init.d/umdns restart >/dev/null 2>&1 || true

exit 0
