# package/mesh-agent/files/etc/uci-defaults/90-mesh-agent-enable
#!/bin/sh

. /usr/lib/mesh/ath11k-nss.sh

/etc/init.d/umdns enable >/dev/null 2>&1 || true

if [ -f /etc/config/mesh_ac ]; then
	uci -q set mesh_agent.main.enabled='0'
	uci -q commit mesh_agent
	/etc/init.d/mesh-agent disable >/dev/null 2>&1 || true
else
	ath11k_nss_set_offload 0 "managed AP first boot" || true
	/usr/sbin/mesh-agent-apply --network-only >/dev/null 2>&1 || true
	/etc/init.d/mesh-agent enable >/dev/null 2>&1 || true
fi

exit 0
