# package/easymesh-agent/files/etc/uci-defaults/90-easymesh-agent-enable
#!/bin/sh

. /usr/lib/easymesh/ath11k-nss.sh

/etc/init.d/umdns enable >/dev/null 2>&1 || true

if [ -f /etc/config/easymesh ]; then
	uci -q set easymesh_agent.main.enabled='0'
	uci -q commit easymesh_agent
	/etc/init.d/easymesh-agent disable >/dev/null 2>&1 || true
	/etc/init.d/easymesh-backhaul enable >/dev/null 2>&1 || true
else
	if ath11k_nss_mesh_support_available; then
		ath11k_nss_set_offload 1 "AP NSS mesh support" || true
	else
		ath11k_nss_set_offload 0 "AP mesh compatibility fallback" || true
	fi
	/usr/sbin/easymesh-agent-apply --network-only >/dev/null 2>&1 || true
	/etc/init.d/easymesh-agent enable >/dev/null 2>&1 || true
	/etc/init.d/easymesh-backhaul enable >/dev/null 2>&1 || true
fi

exit 0
