# package/mesh-ac-local-member/files/etc/uci-defaults/91-mesh-ac-local-member-enable
#!/bin/sh

. /usr/lib/mesh/ath11k-nss.sh

has_wireless_config() {
	[ -f /etc/config/wireless ] || return 1
	uci -q show wireless 2>/dev/null | grep -q '=wifi-device'
}

has_wifi_driver() {
	local phy

	for phy in /sys/class/ieee80211/*; do
		[ -e "$phy" ] && return 0
	done
	command -v iw >/dev/null 2>&1 || return 1
	iw phy >/dev/null 2>&1
}

has_wireless_config || has_wifi_driver || exit 0

if [ "$(uci -q get mesh_ac.main.local_member 2>/dev/null)" = "1" ]; then
	ath11k_nss_set_offload 0 "AC local member first boot" || true
else
	ath11k_nss_set_offload 1 "AC controller mode" || true
fi

exit 0
